# app/vfe-vdpa/vfe-vhostd.service.in
[Unit]
Description=Nvidia Vhost Acceleration Daemon
After=network.target network.service networking.service

[Service]
Type=simple
ExecStartPre=@prefix@/bin/check_pf_reset.sh
ExecStart=@prefix@/bin/vfe-vhostd -v --file-prefix=vfe-vhostd -a 0000:00:00.0 --log-level=.,debug --log-level=lib.vhost.config,info --log-level=pmd.vdpa.virtio,info --log-level=pmd.virtio.dev,info --vfio-vf-token=cdc786f0-59d4-41d9-b554-fed36ff5e89f -- --client
TimeoutSec=1800
LimitNOFILE=200000
LimitCORE=infinity
KillSignal=SIGTERM
Restart=always
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
