# app/virtio-ha/vfe-vhostd-ha.service.in
[Unit]
Description=Nvidia Vhost Acceleration High Availability Daemon
After=network.target network.service networking.service

[Service]
Type=simple
ExecStart=@prefix@/bin/vfe-vhostd-ha
TimeoutSec=1800
LimitNOFILE=200000
LimitCORE=infinity
KillSignal=SIGTERM
Restart=always
ExecStopPost=/usr/bin/bash -c "rm -rf /tmp/pf_resetting"

[Install]
WantedBy=multi-user.target
